# deploy/cuestats.service
[Unit]
Description=CueStats API
After=network.target mongod.service

[Service]
User=REPLACE_USER
WorkingDirectory=REPLACE_PATH/backend
ExecStart=REPLACE_PATH/backend/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
